# massdevserver.service
[Unit]
Description=Mass Development Server
After=network.target

[Service]
User=massdev
Group=massdev
WorkingDirectory=/opt/massdevserver
Environment="PATH=/opt/massdevserver/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/opt/massdevserver"
Environment="PYTHONUNBUFFERED=1"
Environment="FLASK_APP=streamosc.py"
Environment="FLASK_ENV=production"
ExecStart=/opt/massdevserver/venv/bin/python -m flask run --host=0.0.0.0 --port=7401
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target 